# dnsmasq/supervisor.ini
[supervisord]
nodaemon=true

[program:dnscrypt-okturtles]
command=/usr/sbin/dnscrypt-proxy --local-address=127.0.0.1:1053 --resolver-address=23.226.227.93:443 --provider-name=2.dnscrypt-cert.okturtles.com --provider-key=1D85:3953:E34F:AFD0:05F9:4C6F:D1CC:E635:D411:9904:0D48:D19A:5D35:0B6A:7C81:73CB

[program:dnscrypt-eu]
command=/usr/sbin/dnscrypt-proxy --local-address=127.0.0.1:2053 --resolver-address=77.66.84.233:443 --provider-name=2.dnscrypt-cert.resolver2.dnscrypt.eu --provider-key=3748:5585:E3B9:D088:FD25:AD36:B037:01F5:520C:D648:9E9A:DD52:1457:4955:9F0A:9955

[program:dnsmasq]
command=/usr/sbin/dnsmasq --no-daemon --log-queries --no-resolv --no-hosts --conf-file=/etc/dnsmasq.conf
stderr_events_enabled=true
stderr_logfile=/dev/stdout
stderr_logfile_maxbytes=0
