# unbound/supervisor.ini
[supervisord]
nodaemon=true

[program:dnscrypt.eu]
command=/usr/sbin/dnscrypt-proxy --edns-payload-size=4096 --local-address=127.0.0.1:2053 --resolver-address=176.56.237.171:443 --provider-name=2.dnscrypt-cert.resolver1.dnscrypt.eu --provider-key=67C0:0F2C:21C5:5481:45DD:7CB4:6A27:1AF2:EB96:9931:40A3:09B6:2B8D:1653:1185:9C66
user=dnscrypt

[program:fr.dnscrypt.org]
command=/usr/sbin/dnscrypt-proxy --edns-payload-size=4096 --local-address=127.0.0.1:3053 --resolver-address=212.47.228.136 --provider-name=2.dnscrypt-cert.fr.dnscrypt.org --provider-key=E801:B84E:A606:BFB0:BAC0:CE43:445B:B15E:BA64:B02F:A3C4:AA31:AE10:636A:0790:324D
user=dnscrypt

[program:soltysiak.com]
command=/usr/sbin/dnscrypt-proxy --edns-payload-size=4096 --local-address=127.0.0.1:4053 --resolver-address=178.216.201.222:2053 --provider-name=2.dnscrypt-cert.soltysiak.com --provider-key=25C4:E188:2915:4697:8F9C:2BBD:B6A7:AFA4:01ED:A051:0508:5D53:03E7:1928:C066:8F21
user=dnscrypt

[program:okturtles.com]
command=/usr/sbin/dnscrypt-proxy --edns-payload-size=4096 --local-address=127.0.0.1:5053 --resolver-address=23.226.227.93:443 --provider-name=2.dnscrypt-cert.okturtles.com --provider-key=1D85:3953:E34F:AFD0:05F9:4C6F:D1CC:E635:D411:9904:0D48:D19A:5D35:0B6A:7C81:73CB
user=dnscrypt

[program:unbound]
command=/usr/sbin/unbound -d -v -c /etc/unbound/unbound.conf
stderr_events_enabled=true
stderr_logfile=/dev/stdout
stderr_logfile_maxbytes=0
